# ai/openclaw/openclaw-tunnel/examples/supervisord.conf
[supervisord]
nodaemon=true
logfile=/home/appuser/workspace/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:sing-box]
user=root
directory=/home/appuser/workspace
command=sing-box run -c configs/sing-box.tun.json
autostart=true
autorestart=true
startsecs=2
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=logs/sing-box.supervisord.log
stderr_logfile=logs/sing-box.supervisord.error

[program:openclaw]
user=appuser
directory=/home/appuser/workspace
command=/home/appuser/.local/bin/entrypoint.sh openclaw gateway run --port=18789
autostart=true
autorestart=true
startsecs=2
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=logs/openclaw.supervisord.log
stderr_logfile=logs/openclaw.supervisord.error
environment=HOME="/home/appuser",USER="appuser",LOGNAME="appuser"
