# network/hysteria/examples/supervisord.client.conf
[supervisord]
nodaemon=true
logfile=/home/workspace/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:hysteria]
user=appuser
command=/opt/bin/hysteria client -c /home/workspace/configs/hysteria.client.yaml
autostart=true
autorestart=true
startsecs=2
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/workspace/logs/hysteria.log
stderr_logfile=/home/workspace/logs/hysteria.error
