# network/hysteria/examples/supervisord.server.conf
[supervisord]
nodaemon=true
logfile=/home/workspace/logs/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info
stdout_logfile_maxbytes=128MB
stdout_logfile_backups=10
redirect_stderr=true

[unix_http_server]
file=/run/supervisord.sock
# if it's not run as root
# chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:hysteria-0443]
user=root
command=/opt/bin/hysteria server -c /home/workspace/configs/hysteria.server-0443.yaml
autostart=true
autorestart=true
startsecs=2
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/workspace/logs/hysteria-0443.log
#stderr_logfile=/home/workspace/logs/hysteria-0443.log
redirect_stderr=true
stdout_logfile_maxbytes=128MB
stdout_logfile_backups=10

[program:hysteria-8443]
user=root
command=/opt/bin/hysteria server -c /home/workspace/configs/hysteria.server-8443.yaml
autostart=true
autorestart=true
startsecs=2
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/home/workspace/logs/hysteria-8443.log
#stderr_logfile=/home/workspace/logs/hysteria-8443.log
redirect_stderr=true
stdout_logfile_maxbytes=128MB
stdout_logfile_backups=10
