# containers/vpn-socks5/supervisor.ini
# path: /opt/configs/supervisor.ini

# to enable inet_http_server, copy and edit the following lines to /apps/configs/supervisor.ini
;[inet_http_server]    # inet (TCP) server disabled by default
;port=:1098            # (ip_address:port specifier, *:port for all iface)
;username=user         # (default is no username)
;password=123          # (default is no password)

[program:openvpn]
user=root
# --auth-user-pass configs/openvpn.ovpn.auth (account\npassword)
command=openvpn --auth-nocache --config configs/openvpn.ovpn --askpass configs/openvpn.ovpn.pass
autostart=true
autorestart=true
stdout_logfile=logs/openvpn-client.log
stderr_logfile=logs/openvpn-client.error

[program:socat_tcp-8000]
command=socat TCP-LISTEN:9000,fork TCP:192.168.255.1:9000
autostart=true
autorestart=true
stdout_logfile=logs/socat_tcp-9000.log
stderr_logfile=logs/socat_tcp-9000.error

[program:http_server]
command=python3 -m http.server --directory data/site 8000 # --bind 127.0..0.1
autostart=true
autorestart=true
stdout_logfile=logs/http_server.log
stderr_logfile=logs/http_server.error
