# containers/kcptun/examples/supervisor.ini
[program:kcptun-client]
command=kcptun-client -c configs/kcptun-client.json --remoteaddr remote_host:29900 --localaddr :2022
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0

[program:ssh-socks5]
command=ssh -NC -F configs/ssh.conf -E logs/ssh.log -D 127.0.0.1:1080 remote_host
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0
