# container/openvpn/examples/supervisord.nginx.conf
# path: /apps/configs/supervisord.nginx.conf

[unix_http_server]
file=/var/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

;[inet_http_server]    # inet (TCP) server disabled by default
;port=:9002            # (ip_address:port specifier, *:port for all iface)
;username=user         # (default is no username)
;password=123          # (default is no password)

[supervisorctl]
serverurl=unix:///var/run/supervisord.sock
;port=127.0.0.1:9001

[supervisord]
nodaemon=true
user=root
pidfile=/apps/logs/supervisord.pid

# apps
[program:openvpn]
user=root
#command=openvpn --auth-nocache --config configs/client.ovpn --auth-user-pass configs/client.ovpn.auth
command=openvpn --auth-nocache --config configs/client.ovpn --askpass configs/client.ovpn.pass
autostart=true
autorestart=true
stdout_logfile=logs/openvpn-client.log
stderr_logfile=logs/openvpn-client.error

[program:nginx]
user=root
; -c configs/nginx.conf
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=logs/nginx.log
stderr_logfile=logs/nginx.error
