# containers/vpn-socks5/examples/supervisord-full.ini
# path: /opt/configs/supervisor.ini

# to enable inet_http_server, copy and edit the following lines to /apps/configs/supervisor.ini
;[inet_http_server]    # inet (TCP) server disabled by default
;port=:1098            # (ip_address:port specifier, *:port for all iface)
;username=user         # (default is no username)
;password=123          # (default is no password)

[program:openvpn]
autostart=true
autorestart=true
directory=/opt
stdout_logfile=logs/openvpn-client.log
stderr_logfile=logs/openvpn-client.error
user=root
# --auth-user-pass configs/openvpn.ovpn.auth (account\npassword)
command=openvpn --auth-nocache --config configs/openvpn.ovpn --askpass configs/openvpn.ovpn.pass

[program:wireguard]
autostart=true
autorestart=true
#startsecs=10
#startretries=3
#stopwaitsecs=30
#stopsignal=INT
directory=/opt
environment=WG_QUICK_USERSPACE=no
user=root
#command=wg-quick up wg0 # /etc/wireguard/wg0.conf
command=wg-quick up configs/wg0.conf

[program:socat_tcp-9000]
autostart=true
autorestart=true
command=socat TCP-LISTEN:9000,fork TCP:192.168.255.1:9000
stdout_logfile=logs/socat_tcp-9000.log
stderr_logfile=logs/socat_tcp-9000.error

[program:http_server]
autostart=true
autorestart=true
command=python3 -m http.server --directory data/site 8000 # --bind 127.0..0.1
stdout_logfile=logs/http_server.log
stderr_logfile=logs/http_server.error

[program:dante-server]
autostart=true
autorestart=true
command=sockd -f configs/dante-server.conf -N 2
