# container/openvpn/deployments/supervisor.ini
[program:openvpn]
user=root
command=ovpn_run
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0

[program:crond]
user=root
command=crond -l 0 -f -L /apps/logs/crond.log
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0
